# ops/systemd/amigo.service
[Unit]
Description=Amigo Bun Server
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/amigo/backend
Environment=PATH=/home/ubuntu/.bun/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/var/www/amigo/shared/amigo.env
ExecStart=/usr/bin/env bun run dist/server/index.js
Restart=always
RestartSec=3
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
